# supervisor/sphinx.conf
[program:sphinx]
command=/usr/bin/searchd --nodetach -c /etc/sphinx/sphinx.conf
autostart=true
autorestart=true
user=root
priority=100
redirect_stderr=true
